# deploy/systemd/stroma-ai-watcher.service
[Unit]
Description=StromaAI — vLLM Watcher (Dynamic Slurm Burst Scaling)
Documentation=https://github.com/your-org/ai-flux
After=network-online.target ai-flux-vllm.service
Wants=network-online.target
Requires=ai-flux-vllm.service

[Service]
Type=simple
User=aiflux
Group=aiflux
WorkingDirectory=/opt/ai-flux

# Sources all AI_FLUX_* environment variables.
EnvironmentFile=/opt/ai-flux/config.env

# Python is installed by install.sh into /opt/ai-flux/venv.
ExecStart=/opt/ai-flux/venv/bin/python3 /opt/ai-flux/src/vllm_watcher.py

# Restart=always — the Watcher should always be running; failures are
# usually transient (Slurm CLI timeout, Ray connection flap).
Restart=always
RestartSec=10s
TimeoutStartSec=30s
TimeoutStopSec=30s

# Journal logging — view with: journalctl -u ai-flux-watcher -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-flux-watcher

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ai-flux /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
