# deploy/systemd/stroma-ai-vllm.service
[Unit]
Description=StromaAI — vLLM Inference API Server
Documentation=https://github.com/your-org/stroma-ai
After=network-online.target ray-head.service
Wants=network-online.target
Requires=ray-head.service

[Service]
Type=simple
User=stromaai
Group=stromaai
WorkingDirectory=/opt/stroma-ai

# Sources all STROMA_* environment variables.
EnvironmentFile=/opt/stroma-ai/config.env

# vLLM CPU KV cache threads — set before vllm process starts.
# Matches the value baked into the container %environment but can be
# overridden here without rebuilding the image.
Environment=VLLM_CPU_KV_CACHE_THREADS=${STROMA_VLLM_CPU_KV_THREADS}

# Air-gapped: prevent any HuggingFace Hub download attempts at runtime.
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1

# ---------------------------------------------------------------------------
# vLLM serve command
#
# Key flags and their rationale:
#   --quantization awq          — model is AWQ-quantized; reduces VRAM footprint
#   --tensor-parallel-size 1    — one GPU per node; scaling is data-parallel
#   --worker-use-ray            — deploy as a Ray Serve application so Ray
#                                 workers on Slurm nodes contribute capacity
#   --cpu-offload-gb            — overflow KV cache to CPU RAM; critical for
#                                 high concurrency with long contexts on HPC nodes
#   --kv-cache-dtype auto       — vLLM detects optimal dtype for the GPU;
#                                 auto = fp8 on Ada Lovelace/Hopper, fp16 elsewhere
#   --max-model-len             — token window; tune based on available KV budget
#   --max-num-seqs 64           — concurrent sequence slots; tune for expected load
#   --gpu-memory-utilization    — fraction of VRAM used for model + KV cache
#   -O2                         — piecewise CUDA graphs; best throughput on modern
#                                 NVIDIA GPUs (Volta, Ampere, Ada, Hopper).
#                                 NOTE: mutually exclusive with --enforce-eager
#   --host 0.0.0.0              — nginx on the same host proxies; 0.0.0.0 is safe
#   --trust-remote-code         — required by AWQ model tokenizer; review any
#                                 model update against the upstream repo
#
# --block-size is intentionally omitted so vLLM auto-selects the optimal paged
# attention block size for the detected GPU. Add it explicitly (e.g., 16 or 32)
# if you need to pin the value for a specific architecture.
# ---------------------------------------------------------------------------
# Binaries are installed by install.sh into /opt/stroma-ai/venv.
ExecStart=/opt/stroma-ai/venv/bin/vllm serve \
    ${STROMA_MODEL_PATH} \
    --served-model-name ${STROMA_MODEL_NAME} \
    --quantization ${STROMA_VLLM_QUANTIZATION} \
    --tensor-parallel-size 1 \
    --worker-use-ray \
    --cpu-offload-gb ${STROMA_CPU_OFFLOAD_GB} \
    --kv-cache-dtype ${STROMA_KV_CACHE_DTYPE} \
    --max-model-len ${STROMA_MAX_MODEL_LEN} \
    --max-num-seqs ${STROMA_MAX_NUM_SEQS} \
    --gpu-memory-utilization ${STROMA_GPU_MEM_UTIL} \
    --enable-chunked-prefill \
    --host 0.0.0.0 \
    --port ${STROMA_VLLM_PORT} \
    --api-key ${STROMA_API_KEY} \
    --trust-remote-code \
    -O2

# Do NOT use Restart=always — an OOM or CUDA error should surface for review.
Restart=on-failure
RestartSec=30s
TimeoutStartSec=300s
TimeoutStopSec=60s

# Give vLLM time to drain in-flight requests before SIGKILL
KillMode=mixed
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=stroma-ai-vllm

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/stroma-ai /tmp /share
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
