# deploy/systemd/stroma-ai-gateway.service
[Unit]
Description=StromaAI — OIDC Gateway (Authentication Proxy)
Documentation=https://github.com/Moffitt-Cancer-Center/StromaAI
After=network-online.target stroma-ai-vllm.service
Wants=network-online.target
Requires=stroma-ai-vllm.service

[Service]
Type=simple
User=stromaai
Group=stromaai
WorkingDirectory=/opt/stroma-ai

# Sources all STROMA_* and GATEWAY_* environment variables.
EnvironmentFile=/opt/stroma-ai/config.env

# gateway.py launches uvicorn internally on startup (GATEWAY_PORT, default 9000).
# OIDC_DISCOVERY_URL, VLLM_BACKEND_URL, and STROMA_API_KEY must be set in
# config.env for OIDC validation and backend proxying to function.
ExecStart=/opt/stroma-ai/venv/bin/python3 /opt/stroma-ai/src/gateway.py

# Restart=always — transient failures (Keycloak JWKS fetch on startup, network
# blip) should recover automatically.  The lifespan handler tolerates a missing
# IdP at boot and retries on the first request, so we can start without waiting
# for Keycloak to be reachable.
Restart=always
RestartSec=10s
TimeoutStartSec=30s
TimeoutStopSec=30s

# Journal logging — view with: journalctl -u stroma-ai-gateway -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stroma-ai-gateway

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/stroma-ai /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
# ProtectHostname requires systemd >= 243 (RHEL 8 ships 239 — omitted)
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
